# packaging/sdbd_device.service
[Unit]
Description=sdbd
After=default.target

[Service]
Type=forking
PIDFile=/tmp/.sdbd.pid
RemainAfterExit=yes
ExecStart=/usr/sbin/sdbd
