# packaging/sdbd@.socket
[Unit]
Description=sdbd Socket for %I

[Socket]
ListenUSBFunction=%I
Service=sdbd.service
