# packaging/sdbd_device.service
[Unit]
Description=sdbd

[Service]
Type=forking
EnvironmentFile=-/run/tizen-system-env
PIDFile=/tmp/.sdbd.pid
RemainAfterExit=yes
ExecStart=/usr/sbin/sdbd
