# packaging/sdbd_device.service
[Unit]
Description=sdbd
After=syslog.target

[Service]
Type=forking
Environment=DISPLAY=:0
PIDFile=/tmp/.sdbd.pid
RemainAfterExit=yes
ExecStart=/usr/sbin/sdbd

[Install]
WantedBy=multi-user.target
